# deploy/systemd/noah-api.service
[Unit]
Description=Noah ERP API
After=network.target postgresql.service valkey.service

[Service]
Type=simple
WorkingDirectory=/opt/noah-erp/Noah-ERP/apps/api
EnvironmentFile=/opt/noah-erp/Noah-ERP/apps/api/.env
ExecStart=/usr/bin/node dist/main.js
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
